# deploy/beers-crawler.service
# Template unit. scripts/deploy.sh writes a concrete unit with paths from deploy.env.
[Unit]
Description=beers-crawler Untappd API (uvicorn)
After=network.target

[Service]
Type=simple
User=DEPLOY_SERVICE_USER
Group=DEPLOY_SERVICE_GROUP
WorkingDirectory=DEPLOY_ROOT/app
EnvironmentFile=-DEPLOY_ROOT/env
Environment=BEERS_CRAWLER_DB=DEPLOY_ROOT/data/beers.db
Environment=PATH=DEPLOY_ROOT/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=DEPLOY_ROOT/.local/bin/uv run uvicorn beers_crawler.api:app --host 127.0.0.1 --port 8741 --workers 1
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
MemoryMax=512M
StandardOutput=append:DEPLOY_ROOT/logs/api.log
StandardError=append:DEPLOY_ROOT/logs/api.log

[Install]
WantedBy=multi-user.target
